# packages/dcos-check-runner/extra/dcos-checks-poststart.timer
[Unit]
Description=DC/OS Checks Timer: timer for DC/OS Checks service
[Timer]
OnCalendar=*:0/2
AccuracySec=10
