# packages/vault/extra/dcos-vault.service
[Unit]
Description=Vault: DC/OS Default Secret Store Backend

[Service]
User=dcos_vault
PermissionsStartOnly=True
Restart=always
StartLimitInterval=0
RestartSec=5
LimitNOFILE=16384
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=/opt/mesosphere/etc/dcos-vault.env
EnvironmentFile=-/run/dcos/etc/dcos-vault_default.env
ExecStartPre=/opt/mesosphere/bin/bootstrap dcos-vault_default
ExecStart=/opt/mesosphere/bin/vault server -config=/run/dcos/etc/vault.hcl -log-level=trace
