# packages/dcos-signal/extra/dcos-signal.service
[Unit]
Description=Signal: DCOS Telemetry and Tracking Service
[Service]
PermissionsStartOnly=True
Type=simple
Restart=on-failure
RestartSec=20
LimitNOFILE=16384
User=dcos_signal
EnvironmentFile=/opt/mesosphere/environment
EnvironmentFile=-/opt/mesosphere/etc/cfn_signal_metadata
ExecStartPre=/opt/mesosphere/bin/bootstrap dcos-signal
ExecStart=/opt/mesosphere/bin/dcos-signal -segment-key 51ybGTeFEFU1xo6u10XMDrr6kATFyRyh
