# packages/secrets/extra/dcos-secrets.socket
[Unit]
Description=Secrets Service socket: DC/OS Secrets Service socket

[Socket]
ListenStream=/run/dcos/dcos-secrets.sock

[Install]
WantedBy=dcos-secrets.service
